# dell-bios-fan-control/rc.local
#!/bin/sh -e
sh '/usr/local/sbin/dell-bios-fan-control/disable-bios-fan-control.sh'
exit 0
